# debian/cron.d
# Every minute, do nginx-ssl actions
* * * * *	   root		/usr/lib/alternc/update_nginx-ssl.sh
